# protected/extensions/crontab/crontabs/my_crontab
30 22 * * * php /var/www/html/VmiT/yiicmd.php mail
